# config/ss-go.conf
[program:ss-go]
command = /home/blueyi/github/my-utils/bin/ss/shadowsocks-local-linux64
directory = /home/blueyi/github/my-utils/bin/ss
user = root
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/ss-go.log
stderr_logfile = /var/log/supervisor/ss-go_err.log
